# scripts/powerfox-daily.timer
[Unit]
Description=Trigger a daily request to the powerfox-API.

[Timer]
OnCalendar=*-*-* 23:30:00
Persistent=true

[Install]
WantedBy=timers.target