# docs/powerfox-daily.timer
[Unit]
Description=Trigger a daily request to the powerfox-API.

[Timer]
OnCalendar=*-*-* 12:00:00
Persistent=true

[Install]
WantedBy=timers.target